# Configs/runit/service/battery05/run
#!/bin/bash

2>&1

snooze -H'*' -M'*' notify-battery 5
