# Configs/runit/service/dwm/run
#!/bin/sh

#export TZ=America/Denver
exec 2>&1
cd ~
exec /home/richard/git/dwm/dwm
