# Configs/runit/service/clipcatd/run
#!/bin/bash

2>&1

exec clipcatd --no-daemon
