# Configs/runit/service/slstatus/run
#!/bin/sh

#export TZ=America/Denver
exec 2>&1
exec slstatus
